# data/io.github.kolunmi.bazaar.service.in
[Unit]
Description=Bazaar Service
PartOf=graphical-session.target

[Service]
Type=dbus
BusName=io.github.kolunmi.bazaar
ExecStart=@bindir@/bazaar service
