# data/io.github.kolunmi.Bazaar.service.in
[Unit]
Description=Bazaar Service
PartOf=graphical-session.target

[Service]
Type=dbus
BusName=io.github.kolunmi.Bazaar
ExecStart=@bindir@/bazaar service
