# data/io.github.kolunmi.bazaar.service.in
[Unit]
Description=Bazaar Service
ConditionUser=!@system
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=@bindir@/bazaar service

[Install]
WantedBy=multi-user.target
